# code/raspberry-pi-hardware-interactions-part-2-rotary-encoder/volume-dial.service
[Unit]
Description=Volume Dial
After=sysinit.target

[Service]
ExecStart=/usr/bin/python /opt/volume_dial.py

[Install]
WantedBy=multi-user.target
